# supervisord.conf
[supervisord]
nodaemon=true

[program:sinatra]
command=rackup
autorestart=true
